# overlay/etc/systemd/system/gitlab-runsvdir.path
[Unit]
Description="Replacing incron with systemd."

[Path]
DirectoryNotEmpty=/
Unit=gitlab-runsvdir.service

[Install]
WantedBy=basic.target
